# etc/systemd/system/tlp.service.d/override.conf
[Service]
Restart=
Restart=always
RestartSec=
RestartSec=1s
