# etc/systemd/system/rpm-ostreed-automatic.service.d/override.conf
[Service]
ExecStart=
ExecStart=/etc/redneck-bin/update.fish
User=root
Nice=19
IOSchedulingClass=best-effort
IOSchedulingPriority=7
RestartSec=1s
RestrictRealtime=true
