# LXroot/etc/systemd/system/autopgrade.service
[Unit]
Description=Automatic upgrader
After=basic.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/topgrade --skip-notify -y
IOSchedulingPriority=7
Nice=19

[Install]
WantedBy=multi-user.target
