# LXroot/etc/systemd/system/boinc-client.service
[Unit]
Description=BOINC Client for contributing idle device resources
After=basic.target

[Service]
Type=simple
User=root
ReadWritePaths=-/var/lib/boinc -/etc/boinc-client -/tmp
WorkingDirectory=/var/lib/boinc
ExecStart=/bin/boinc
ExecStopPost=/bin/rm -f lockfile
IOSchedulingPriority=7
Restart=always
RestartSec=1
Nice=19

[Install]
WantedBy=multi-user.target
