# LXroot/etc/systemd/system/fyn-refyne.service
[Unit]
Description=Timed refinement of your system. May be risky in production environments.
After=network.target

[Service]
User=root
ExecStart=sudo curl https://raw.githubusercontent.com/MrGrappleMan/Fynelium/main/start.bash | sudo bash

[Install]
WantedBy=multi-user.target
