# LXroot/etc/systemd/system/swappity.timer
[Unit]
Description=Timer for Dynamic Swap and Memory Adjuster

[Timer]
OnBootSec=30s
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target