# LXroot/etc/systemd/system/zram.service
[Unit]
Description=ZRAM Starter
After=basic.target

[Service]
Type=oneshot
ExecStart=/etc/silverbluebins/zram.fish
Nice=-15
User=root

[Install]
WantedBy=multi-user.target
