# etc/systemd/system/.service
[Unit]
Description=Package manager updater
After=network-online.target
ConditionPathExists=/run/ostree-booted
Wants=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/fish /etc/redneck-bin/pkgm.fish
IOSchedulingPriority=7
Nice=19

[Install]
WantedBy=multi-user.target
