# etc/systemd/system/fn-mcbe-srv.service
[Service]
Type=simple
WorkingDirectory=/opt/minecraft-bedrock
ExecStart=/usr/bin/tmux new-session -s bedrock -d '/opt/mcbe-server/bedrock_server'
ExecStop=/usr/bin/tmux kill-session -t bedrock
RemainAfterExit=yes
StandardInput=tty
StandardOutput=journal
Nice=19

[Install]
