# etc/systemd/system/boinc-client.service.d/override.conf
[Service]
Restart=always
RestartSec=1s
Nice=
Nice=19
IOSchedulingPriority=7
IOSchedulingClass=
IOSchedulingClass=best-effort
RestrictRealtime=true
