# etc/systemd/system/boinc-client.service.d/override.conf
[Service]
Nice=
Nice=19
IOSchedulingClass=
IOSchedulingClass=best-effort
IOSchedulingPriority=7
