# etc/systemd/system/boinc-client.service.d/override.conf
[Service]
Nice=
Nice=19
IOSchedulingClass=
IOSchedulingClass=best-effort
IOSchedulingPriority=
IOSchedulingPriority=7
Restart=
Restart=always
RestartSec=
RestartSec=1s
